# getmail.cron
*/5 * * * * root getmail -g /etc/getmail